# server/wg-mgmt.service
[Unit]
Description=WireGuard Management Daemon
After=network-online.target wg-quick@wg0.service
Wants=network-online.target wg-quick@wg0.service
Requires=wg-quick@wg0.service

[Service]
Type=simple
User=root
# Template — setup-server.sh generates the real unit with actual paths.
# Do not deploy this file directly.
ExecStart=/usr/local/bin/wg-mgmt-daemon --config=/home/yequdesu/WG-manager/config.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
